# printer-test-page.timer
# systemd user timer that runs the test page once a week
# Install alongside printer-test-page.service (see that file for instructions).

[Unit]
Description=Monthly printer test page to keep heads from clogging

[Timer]
# See man systemd.time
# Run on the 15th of every month at 8am
OnCalendar=*-*-15 08:00:00
Persistent=true

[Install]
WantedBy=timers.target

